# conf/playout.supervised.conf
[group:playout]
programs=liquidsoap,pypo
priority=100

[program:liquidsoap]
environment=HOME='/home/playout/'
directory=/home/playout/pypo/src/pypo/liquidsoap_scripts
command=/usr/local/bin/liquidsoap --verbose -f ls_script.liq
user=playout
autostart=true
autorestart=true
redirect_stderr=True
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stdout_logfile=/var/log/playout/supervisor-liquidsoap.log


[program:pypo]
environment=HOME='/home/playout/'
directory=/home/playout/pypo/src/pypo/
command=/home/playout/env/pypo/bin/python pypo.py
user=playout
autostart=true
autorestart=true
redirect_stderr=True
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stdout_logfile=/var/log/playout/supervisor-pypo.log
