# conf/playout.supervised.conf
[group:playout]
programs=liquidsoap,pypo
priority=100

[program:liquidsoap]
directory=/home/playout/pypo/src/pypo/liquidsoap_scripts
command=/usr/local/bin/liquidsoap --verbose -f ls_script.liq
user=playout
autostart=true
autorestart=true
redirect_stderr=True
environment=HOME='/home/playout/'

[program:pypo]
directory=/home/playout/pypo/src/pypo/
command=/home/playout/env/pypo/bin/python pypo.py
user=playout
autostart=true
autorestart=true
redirect_stderr=True
environment=HOME='/home/playout/'
